# .fetchmail.cron
*/5 * * * * /usr/bin/fetchmail --pidfile /home/vmiklos/fetchmail.pid
